# root/etc/services.d/socat/run
#!/usr/bin/with-contenv sh

umask 002

s6-echo "Forwarding $LISTEN_PORT -> target_app:$FORWARD_PORT"

exec socat tcp-listen:${LISTEN_PORT},fork tcp-connect:target_app:${FORWARD_PORT}
